# roles/common/templates/etc/systemd/system/pipy.service
# This file is automatically generated by Ansible.
# {{ ansible_managed }}
[Unit]
Description=PiPy Web Service

[Service]
Restart=always
ExecStart={{ python_bin }}/twistd --nodaemon \
          --uid {{ pipy_user }} --gid {{ pipy_user }} web \
          --port le:{{ pipy_certs_dir }}:tcp:443

[Install]
WantedBy=multi-user.target
